# profile.d/path.sh
#!/usr/bin/bash.exe
PATH=/c/Users/Public/OneDrive/Documents/Projects/Python/Packages/PyClasses/scripts:${PATH}