# profile.d/PyClasses.sh
#!/usr/bin/bash.exe
echo "Load PyClasses profile!"